# build_install/files/10-autolvm.rules
ACTION=="add", KERNEL=="vd*[!0-9]", RUN+="/usr/sbin/autolvm REPLACE THIS"
